# scripts/wled-backend.service
# Systemd service file for WLED Control Interface Backend
# Copy this file to /etc/systemd/system/wled-backend.service
# Then enable and start with:
#   sudo systemctl daemon-reload
#   sudo systemctl enable wled-backend
#   sudo systemctl start wled-backend

[Unit]
Description=WLED Control Interface Backend API
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/ledcontrol/backend
Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=/var/www/ledcontrol/backend/.env
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wled-backend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
